# infra/supervisord.conf
; Supervisor Configuration for Coding Engine Container
; Manages all processes: VNC, Control API, and Coding Engine

[supervisord]
nodaemon=true
logfile=/data/logs/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/data/logs
user=root

; ===== VNC Stack =====

[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 %(ENV_DISPLAY_WIDTH)sx%(ENV_DISPLAY_HEIGHT)sx%(ENV_DISPLAY_DEPTH)s -ac +extension GLX +render -noreset
autorestart=true
priority=100
stdout_logfile=/data/logs/xvfb.log
stderr_logfile=/data/logs/xvfb.err
startsecs=2

[program:fluxbox]
command=/usr/bin/fluxbox -display :99
autorestart=true
priority=200
environment=DISPLAY=":99"
stdout_logfile=/data/logs/fluxbox.log
stderr_logfile=/data/logs/fluxbox.err
startsecs=3
depends_on=xvfb

[program:x11vnc]
command=/usr/bin/x11vnc -display :99 -forever -shared -rfbport %(ENV_VNC_PORT)s -nopw -xkb -noxrecord -noxfixes -noxdamage -permitfiletransfer
autorestart=true
priority=300
environment=DISPLAY=":99"
stdout_logfile=/data/logs/x11vnc.log
stderr_logfile=/data/logs/x11vnc.err
startsecs=5
depends_on=xvfb

[program:novnc]
command=/opt/novnc/utils/novnc_proxy --vnc localhost:%(ENV_VNC_PORT)s --listen %(ENV_NOVNC_PORT)s --web /opt/novnc
autorestart=true
priority=400
stdout_logfile=/data/logs/novnc.log
stderr_logfile=/data/logs/novnc.err
startsecs=5
depends_on=x11vnc

; ===== Control API Server =====

[program:control-api]
command=/opt/venv/bin/uvicorn control_server.server:app --host %(ENV_CONTROL_HOST)s --port %(ENV_CONTROL_PORT)s --log-level info
directory=/app
autorestart=true
priority=500
environment=
    PATH="/opt/venv/bin:%(ENV_PATH)s",
    DISPLAY=":99",
    OUTPUT_DIR="%(ENV_OUTPUT_DIR)s",
    REQUIREMENTS_DIR="%(ENV_REQUIREMENTS_DIR)s",
    LOG_DIR="%(ENV_LOG_DIR)s",
    NO_TIMEOUT="%(ENV_NO_TIMEOUT)s",
    MAX_ITERATIONS="%(ENV_MAX_ITERATIONS)s"
stdout_logfile=/data/logs/control-api.log
stderr_logfile=/data/logs/control-api.err
startsecs=5

; ===== Terminal for CLI Preview =====

[program:terminal]
command=/usr/bin/xfce4-terminal --display=:99 --geometry=120x40+50+50 --title="Coding Engine Terminal"
autorestart=false
autostart=false
priority=600
environment=DISPLAY=":99"
stdout_logfile=/data/logs/terminal.log
stderr_logfile=/data/logs/terminal.err

; ===== Browser for Web Preview =====

[program:browser]
command=/usr/bin/chromium-browser --display=:99 --no-sandbox --disable-gpu --disable-software-rasterizer --window-size=1280,800 --window-position=100,100 http://localhost:%(ENV_PREVIEW_PORT)s
autorestart=false
autostart=false
priority=700
environment=DISPLAY=":99"
stdout_logfile=/data/logs/browser.log
stderr_logfile=/data/logs/browser.err

; ===== Process Groups =====

[group:vnc]
programs=xvfb,fluxbox,x11vnc,novnc
priority=100

[group:services]
programs=control-api
priority=200

[group:preview]
programs=terminal,browser
priority=300

; ===== Event Listener for Logging =====

[eventlistener:process-monitor]
command=/app/scripts/process-monitor.sh
events=PROCESS_STATE_EXITED,PROCESS_STATE_FATAL
buffer_size=100
stdout_logfile=/data/logs/process-monitor.log
stderr_logfile=/data/logs/process-monitor.err